# server/cherrypy.conf
[supervisord]

[program:deepspeech-server]
command=gunicorn wsgi --config /deepspeech/server/gunicorn.conf --chdir=/deepspeech/server
